# sh/initd-header.sh
#!/bin/sh

### BEGIN INIT INFO
# Provides:          karuta-server
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: karuta-server
# Description:       An extensive websocket server
### END INIT INFO
